# systemd/vacuum-replays.service
[Unit]
Description=Remove old server replays

[Service]
Type=oneshot
ExecStart=/root/vacuum-replays.sh
