# systemd/vacuum-replays.timer
[Unit]
Description=Remove old server replays

[Timer]
Unit=vacuum-replays.service
OnCalendar=*-*-* 07:00

[Install]
WantedBy=multi-user.target
